# packaging/systemd/onibi.service
[Unit]
Description=Onibi daemon
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/onibi --headless --auto-transports
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info,onibi=info

[Install]
WantedBy=default.target
